# snapp/vdirsyncer/vdirsyncer.service
[Unit]
Description=Sync calendar with vdirsyncer

[Service]
ExecStart=vdirsyncer sync
