# guake/guake.service
[Unit]
Description=Start guake as a Service

[Service]
ExecStart=/bin/bash guake
