# guake/guake.service
[Unit]
Description=Start guake as a Service

[Service]
ExecStart=guake

[Install]
WantedBy=graphical.target
