# guake/quake.service
[Unit]
Description=Start guake

[Service]
ExecStart=/bin/bash guake

[Install]
WantedBy=graphical.target
