# i3/services/feh.service
[Unit]
Description=Change background

[Service]
ExecStart=feh --bg-fill --randomize %h/Pictures/GoSiMac/
