# sway/swaybg/swaybg.service
[Unit]
Description=Change background

[Service]
ExecStart=%h/.config/swaybg/gosimac.sh
