# tmpfiles/systemd-tmpfiles-clean.timer.d/override.conf
[Unit]
Description=Hourly Cleanup of User's Temporary Directories

[Timer]
OnStartupSec=5min
OnUnitActiveSec=1h
