# i3/nitrogen/nitrogen.service
 [Unit]
Description=Change background

[Service]
ExecStart=nitrogen --set-scaled --random %h/Pictures/GoSiMac/
