# i3/nitrogen/nitrogen.service
 [Unit]
Description=Change background

[Service]
ExecStart=nitroge --random %h/Pictures/GoSiMac/
