# i3/nitrogen/nitrogen.service
[Unit]
Description=Change background

[Service]
ExecStart=nitrogen --set-zoom-fill --random %h/Pictures/GoSiMac/
